# animal-humane/deployment/systemd/animal-humane-scheduler.service
# Systemd service file for running the scheduler as a system service
# Copy this to /etc/systemd/system/ and adjust paths as needed

[Unit]
Description=Animal Humane Background Scheduler
After=network.target elasticsearch.service
Wants=elasticsearch.service

[Service]
Type=simple
User=your-username
Group=your-group
WorkingDirectory=/path/to/your/animal-humane
Environment=PYTHONPATH=/path/to/your/animal-humane
ExecStart=/usr/bin/python3 /path/to/your/animal-humane/scheduler/background_scheduler.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Create necessary directories
ExecStartPre=/bin/mkdir -p /path/to/your/animal-humane/logs
ExecStartPre=/bin/mkdir -p /path/to/your/animal-humane/diff_reports

[Install]
WantedBy=multi-user.target